# ultron.service
[Unit]
Description=ULTRON Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/agent/ULTRON/ULTRON
Environment="PATH=/agent/.local/share/fnm/node-versions/v24.11.1/installation/bin:/agent/ULTRON/ULTRON/node_modules/.bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/agent/ULTRON/ULTRON/.env
ExecStart=/agent/ULTRON/ULTRON/node_modules/.bin/tsx /agent/ULTRON/ULTRON/src/index.ts
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
